# packages/smart-home/wyoming/openwakeword/rootfs/etc/s6-overlay/s6-rc.d/openwakeword/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start openWakeWord service
# ==============================================================================
flags=()

if [[ "${OPENWAKEWORD_DEBUG}" == "true" ]]; then
    flags+=('--debug')
fi

exec python3 -m wyoming_openwakeword \
    --uri "tcp://0.0.0.0:${OPENWAKEWORD_PORT}" \
    --preload-model "${OPENWAKEWORD_PRELOAD_MODEL}" \
    --custom-model-dir "${OPENWAKEWORD_CUSTOM_MODEL_DIR}" \
    --threshold "${OPENWAKEWORD_THRESHOLD}" \
    --trigger-level "${OPENWAKEWORD_TRIGGER_LEVEL}" ${flags[@]}
